# monitor/check_srcds_monitor@.service
[Unit]
Description=SRCDS monitor for instance %i
After=network.target

[Service]
Type=simple
# Per-instance EnvironmentFile (create this at ~/.config/systemd/user/check_srcds_monitor@%i.env)
# The env file should at minimum define RESTART_UNIT (SERVER_HOST and PORT are optional).
EnvironmentFile=%h/.config/systemd/user/check_srcds_monitor@%i.env
ExecStart=/usr/bin/env python3 /home/youruser/check_srcds_restart.py
# Use SIGINT as the stop signal so the script receives the same signal it expects for graceful shutdown (Ctrl+C)
KillSignal=SIGINT
TimeoutStopSec=20
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
