# scripts/sourcebans-update.service
[Unit]
Description=Update the SourceBans web components
After=network.target

[Service]
ExecStart=/opt/deploy_sourcebans.sh
